# ics-part/supervisord.conf
#=========================================================================================
# Alumno: Ibai Castañon Osorio
# Email: icasto01@estudiantes.unileon.es
#
# Descripcion: Fichero de configuracion de supervisor, aqui definimos todas los elementos
#              que queremos iniciar al arranque del contenedor
#=========================================================================================
[supervisord]
nodaemon = true

[program:sshd]
startretries=3
command=/usr/sbin/sshd -D

[program:start-honeyd]
startretries=3
command=/bin/bash -c /start-honeyd.sh