# ics-part/supervisord.conf
#=========================================================================================
# Alumno: Ibai Castañon Osorio
# Email: icasto01@estudiantes.unileon.es
#
# Supervisor configuration file, here we define all the elements that we want to start 
# at container startup
#=========================================================================================
[supervisord]
nodaemon = true

[program:sshd]
startretries=3
command=/usr/sbin/sshd -D

[program:start-honeyd]
startretries=3
command=/bin/bash -c /start-honeyd.sh